# gost/files/gost.init
#!/bin/sh /etc/rc.common
# Created By ImmortalWrt
# https://github.com/project-openwrt

START=90
STOP=10

enable="$(uci get gost.@gost[0].enable)"
run_command="$(uci get gost.@gost[0].run_command)"

start()
{
	stop

	if [ "$enable" -eq "1" ]; then
        /usr/bin/gost ${run_command} >/dev/null 2>&1 &
  fi
}

stop()
{
	killall -9 "gost" > "/dev/null" 2>&1
}
